# Data/crontab
*/5  *   *   *   *    curl -s "http://127.0.0.1:9117/jsondb/save"

*/15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/rutor/parse"
*    */4 *   *   *    curl -s "http://127.0.0.1:9117/cron/rutor/UpdateTasksParse"
*/5  *   *   *   *    curl -s "http://127.0.0.1:9117/cron/rutor/ParseAllTask"

*/15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/rutracker/parse"
*    */4 *   *   *    curl -s "http://127.0.0.1:9117/cron/rutracker/UpdateTasksParse"
*    */4 *   *   *    curl -s "http://127.0.0.1:9117/cron/rutracker/ParseAllTask"
*/5  *   *   *   *    curl -s "http://127.0.0.1:9117/cron/rutracker/parseMagnet"

*/15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/kinozal/parse"
*    */4 *   *   *    curl -s "http://127.0.0.1:9117/cron/kinozal/UpdateTasksParse"
*/5  *   *   *   *    curl -s "http://127.0.0.1:9117/cron/kinozal/ParseAllTask"
*    */4 *   *   *    curl -s "http://127.0.0.1:9117/cron/kinozal/parseMagnet"

*/15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/nnmclub/parse"
*    */4 *   *   *    curl -s "http://127.0.0.1:9117/cron/nnmclub/UpdateTasksParse"
*/5  *   *   *   *    curl -s "http://127.0.0.1:9117/cron/nnmclub/ParseAllTask"

*/15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/selezen/parse"
*    */4 *   *   *    curl -s "http://127.0.0.1:9117/cron/selezen/UpdateTasksParse"
*/5  *   *   *   *    curl -s "http://127.0.0.1:9117/cron/selezen/ParseAllTask"

*/30 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/megapeer/parse"
*    */4 *   *   *    curl -s "http://127.0.0.1:9117/cron/megapeer/UpdateTasksParse"
*/15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/megapeer/ParseAllTask"

*/15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/torrentby/parse"
*    */4 *   *   *    curl -s "http://127.0.0.1:9117/cron/torrentby/UpdateTasksParse"
*/5  *   *   *   *    curl -s "http://127.0.0.1:9117/cron/torrentby/ParseAllTask"

*/40 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/toloka/parse"
*    */4 *   *   *    curl -s "http://127.0.0.1:9117/cron/toloka/parseMagnet"

*/15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/lostfilm/parse"

# Bitru: старый парсинг по HTML (browse.php)
# */15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/bitru/parse"
# *    */4 *   *   *    curl -s "http://127.0.0.1:9117/cron/bitru/UpdateTasksParse"
# */5  *   *   *   *    curl -s "http://127.0.0.1:9117/cron/bitru/ParseAllTask"
# *    */4 *   *   *    curl -s "http://127.0.0.1:9117/cron/bitru/parseMagnet"

# Bitru: официальный API (5 req/sec), movie+serial
*/20 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/bitruapi/parse"
# Разово подтянуть торренты после даты: curl -s "http://127.0.0.1:9117/cron/bitruapi/ParseFromDate?lastnewtor=23.10.2025"

# */15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/animelayer/parse"
# */15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/animedia/parse"

# RIP (retired) — anilibria, anifilm, hdrezka: disabled by default via disable_trackers
# */15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/anilibria/parse"
# */15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/anifilm/parse"
# */15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/hdrezka/parse"

# */15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/anidub/parse"


# */15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/baibako/parse"

# */15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/hamsterstudio/parse"
