# Data/crontab
# Save JSON DB
*/5  *   *   *   *    curl -s "http://127.0.0.1:9117/jsondb/save"

# Rutor
*/15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/rutor/parse"
*    */4 *   *   *    curl -s "http://127.0.0.1:9117/cron/rutor/UpdateTasksParse"
*/5  *   *   *   *    curl -s "http://127.0.0.1:9117/cron/rutor/ParseAllTask"

# Rutracker
*/15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/rutracker/parse"
*    */4 *   *   *    curl -s "http://127.0.0.1:9117/cron/rutracker/UpdateTasksParse"
*    */4 *   *   *    curl -s "http://127.0.0.1:9117/cron/rutracker/ParseAllTask"

# Kinozal
*/15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/kinozal/parse"
*    */4 *   *   *    curl -s "http://127.0.0.1:9117/cron/kinozal/UpdateTasksParse"
*/5  *   *   *   *    curl -s "http://127.0.0.1:9117/cron/kinozal/ParseAllTask"

# Nnmclub
*/15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/nnmclub/parse"
*    */4 *   *   *    curl -s "http://127.0.0.1:9117/cron/nnmclub/UpdateTasksParse"
*/5  *   *   *   *    curl -s "http://127.0.0.1:9117/cron/nnmclub/ParseAllTask"

# Selezen
*/15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/selezen/parse"
# curl "http://127.0.0.1:9117/cron/selezen/parse?parseFrom=1&parseTo=5"

# Megapeer
0   *   *   *   *    curl -s "http://127.0.0.1:9117/cron/megapeer/parse"
5   */4 *   *   *    curl -s "http://127.0.0.1:9117/cron/megapeer/UpdateTasksParse"
30  *   *   *   *    curl -s "http://127.0.0.1:9117/cron/megapeer/ParseAllTask"

# Torrentby
*/15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/torrentby/parse"
*    */4 *   *   *    curl -s "http://127.0.0.1:9117/cron/torrentby/UpdateTasksParse"
*/5  *   *   *   *    curl -s "http://127.0.0.1:9117/cron/torrentby/ParseAllTask"

# Toloka
*/40 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/toloka/parse"

# Mazepa
*/40 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/mazepa/parse"

# Lostfilm
*/15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/lostfilm/parse"
# curl "http://127.0.0.1:9117/cron/lostfilm/ParsePages?pageFrom=1&pageTo=10"

# BitruApi: официальный API (5 req/sec), movie+serial
*/20 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/bitruapi/parse"

# Knaben: API v1, TV + Movies
*/20 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/knaben/parse"
# pages=3 | query=Breaking+Bad | hours=24 | orderBy=seeders | categories=2001000,3001000

# Animelayer
*/15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/animelayer/parse"
# curl "http://127.0.0.1:9117/cron/animelayer/parse?parseFrom=1&parseTo=5"

# Anidub
*/15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/anidub/parse"
# curl "http://127.0.0.1:9117/cron/anidub/parse?parseTo=5"
# curl "http://127.0.0.1:9117/cron/anidub/parse?parseFrom=1&parseTo=5"

# Aniliberty
*/15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/aniliberty/parse"
# curl "http://127.0.0.1:9117/cron/aniliberty/parse?parseTo=5"
# curl "http://127.0.0.1:9117/cron/aniliberty/parse?parseFrom=1&parseTo=5"
# https://aniliberty.top/api/v1/anime/torrents?limit=50&page=60

# Baibako
*/15 *   *   *   *    curl -s "http://127.0.0.1:9117/cron/baibako/parse"
# curl "http://127.0.0.1:9117/cron/baibako/parse?parseFrom=0&parseTo=10"
# curl "http://127.0.0.1:9117/cron/baibako/parse?parseFrom=5"  # Parse single page 5
